# config/com.ben-digest.nullclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ben-digest.nullclaw</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/nullclaw</string>
        <string>gateway</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/ben-digest-nullclaw.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ben-digest-nullclaw.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ben</string>
    </dict>
</dict>
</plist>
